# bot/systemd/endcord-bot.service.example
[Unit]
Description=Endcord companion bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/endcord/bot
ExecStart=/opt/endcord/bot/scripts/run-local-runtime.sh /opt/endcord/bot/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
